# stage-ble/00-stage-ble/files/run_boot_script.service
[Unit]
Description=Run custom boot commands
After=network.target

[Service]
Type=oneshot
ExecStart=/home/pi/boot_script.sh
RemainAfterExit=yes
User=root

[Install]
WantedBy=multi-user.target